# ops/systemd/vigil-cleanup.service
[Unit]
Description=Krypton Vigil cleanup (purge expired screenshots/recordings)
After=krypton-vigil.service

[Service]
Type=oneshot
WorkingDirectory=/opt/krypton-vigil
ExecStart=/opt/krypton-vigil/.venv/bin/python -m app.scripts.cleanup
User=root
Group=root
